# hl7-standards-dev-dashboard.service
[Unit]
Description=HL7 Standards Development Dashboard
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/code/hl7-standards-dev-dashboard
ExecStart=/home/ubuntu/.venv/bin/python run.py
Restart=always
RestartSec=3
Environment=PATH=/home/ubuntu/.venv/bin:/usr/bin:/usr/local/bin

# Logging
StandardOutput=append:/var/log/hl7-standards-dev-dashboard/app.log
StandardError=append:/var/log/hl7-standards-dev-dashboard/error.log

[Install]
WantedBy=multi-user.target
